# dev-util/drone/files/droned-0.2.1.init
#!/sbin/runscript

start() {
  ebegin "Starting droned"
  start-stop-daemon --start --exec /usr/sbin/droned -- --port=:80
  eend $?
}

stop() {
  ebegin "Starting droned"
  start-stop-daemon --stop --exec /usr/sbin/droned
  eend $?
}
